# config/systemd/os-system-agent-daily.timer.example
# OS_SYSTEM_AGENT — daily ETL report timer (systemd USER timer).
# Pairs with os-system-agent-daily.service. Install steps are in that file.

[Unit]
Description=Run the OS_SYSTEM_AGENT daily ETL report

[Timer]
# Two digests a day: morning (just after the ETL window) and evening.
# Multiple OnCalendar lines are OR'd — the report fires at each time.
OnCalendar=*-*-* 08:30:00
OnCalendar=*-*-* 18:00:00
# If the machine was asleep/off at a fire time, run as soon as it is back.
Persistent=true
# Spread load / avoid a thundering herd — random 0..120s extra delay.
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
